# services/nightly-valuation.service
# /etc/systemd/system/nightly-valuation.service
[Unit]
Description=Nightly Inventory Valuation Snapshot
Documentation=https://github.com/sundialwire/greenlight
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=welch
Group=welch
WorkingDirectory=/home/welch/projects/sundial_greenlight
Environment="PATH=/home/welch/projects/sundial_greenlight/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/welch/projects/sundial_greenlight/venv/bin/python util/nightly_valuation.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nightly-valuation
